# service/claude-code-proxy.service
[Unit]
Description=Claude Code GLM-5.1 Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=PROXY_DIR_PLACEHOLDER
ExecStart=VENV_PYTHON_PLACEHOLDER -u PROXY_PY_PLACEHOLDER
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/tmp/claude-proxy.log
StandardError=append:/tmp/claude-proxy.log

[Install]
WantedBy=default.target
